# config/systemd/user/mpd.socket.d/override.conf
[Socket]
# Stock unit binds all interfaces, exposing MPD to the LAN. Empty
# ListenStream= clears the inherited list before replacing it below.
ListenStream=
ListenStream=%t/mpd/socket
ListenStream=127.0.0.1:6600
